# deploy/systemd/laravel-schedule.timer
[Unit]
Description=Run Laravel scheduler every minute

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min
AccuracySec=1s

[Install]
WantedBy=timers.target
